# renderer/smc-renderer.service
# Systemd unit for the sun-moon-clock renderer.
#
# Install:
#   sudo cp smc-renderer.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now smc-renderer.service
#
# Logs:
#   journalctl -u smc-renderer -f

[Unit]
Description=Sun Moon Clock renderer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Edit User and WorkingDirectory to match your deployment.
User=smc
WorkingDirectory=/opt/sunMoonClock/renderer

# Port the HTTP server listens on, and the upstream clock URL it
# screenshots. SMC_URL is the public address of the deployed web clock.
Environment=PORT=3000
Environment=SMC_URL=http://localhost:8000

ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
